# source/client/clientLights/client_lights.service
[Unit]
Description=client light service

[Service]
Type=simple
Restart=always
RestartSec=1
User=student
ExecStart=/usr/bin/client-light
